# community/netbook/skel/.config/autostart/xfce4-clipman-plugin-autostart.desktop
[Desktop Entry]
Hidden=false
TryExec=xfce4-clipman
Exec=xfce4-clipman
