# sonar/mate/skel/.config/autostart/sonar_desktop_items.desktop
[Desktop Entry]
Type=Application
Name=Sonar Desktop Items
Exec=/usr/bin/sonar_desktop_items
X-GNOME-Autostart-enabled=true
NoDisplay=false
