# community/sway/etc/systemd/user/way-displays.service
[Unit]
Description=Auto Manages Your Wayland Displays
Documentation=https://github.com/alex-courtis/way-displays
PartOf=wayland-session.target
After=wayland-session.target

[Service]
NonBlocking=true
ExecStart=way-displays --log-threshold debug
ExecReload=kill -SIGUSR2 $MAINPID
Restart=on-failure
StandardOutput=file:/tmp/way-displays.log
StandardError=file:/tmp/way-displays.log

[Install]
WantedBy=wayland-session.target
