# community/sway/etc/systemd/user/way-displays.service
[Unit]
Description=Auto Manages Your Wayland Displays
Documentation=https://github.com/alex-courtis/way-displays
PartOf=wayland-session.target
After=wayland-session.target

[Service]
NonBlocking=true
ExecStart=way-displays > /tmp/way-displays.${XDG_VTNR}.${USER}.log 2>&1
ExecReload=kill -SIGUSR2 $MAINPID
Restart=on-failure

[Install]
WantedBy=wayland-session.target
