# community/lxqt/skel/.config/autostart/claws-mail.desktop
[Desktop Entry]
Exec=claws-mail
Hidden=true
Name=claws-mail
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
