# community/lxqt/skel/.config/autostart/compton.desktop
[Desktop Entry]
Exec=compton -b
Hidden=true
Name=compton
Type=Application
Version=1.0
