# community/lxqt/skel/.config/autostart/live-session.desktop
[Desktop Entry]
Exec=desktop-items
Name=desktop-items
Comment=Live session desktop items
Type=Application
