# community/lxqt/skel/.config/autostart/msm_notifier.desktop
[Desktop Entry]
Exec=msm_notifier
Name=MSM Notifier
StartupNotify=false
Terminal=false
Type=Application
