# community/lxqt/skel/.config/autostart/network-manager-applet.desktop
[Desktop Entry]
Exec=/usr/bin/nm-applet
Name=network-manager-applet
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
