# community/lxqt/skel/.config/autostart/redshiftconf.desktop
[Desktop Entry]
Exec=/usr/bin/redshiftconf
Name=redshiftconf
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
Hidden=true

