# shared/live/skel/.config/autostart/live-session.desktop
[Desktop Entry]
Exec=live-session
Name=live-session
Comment=Live session desktop items
Type=Application
Version=0.9.9

