# community/gnome/skel/.config/autostart/create-template.desktop
[Desktop Entry]
Exec=create-template
Name=create-template
Comment=Create Empty File for Nautilus "New Document" Menu
Type=Application
